# roles/shell/templates/etc/profile.d/custom_prompt.sh.j2
# /etc/profile.d/custom_prompt.sh
# Custom shell prompt for all users
# This file is managed by Ansible (arillso.system.ui)

# Set a universal, informative PS1 prompt
export PS1="\[\e[0;32m\]\u@\h \[\e[0;34m\]\w\[\e[0m\]$ "
# Optionally, add more prompt customization as needed
