# roles/tuning/templates/etc/systemd/system/disable-thp.service.j2
# Transparent Huge Pages configuration service
# Generated by Ansible system tuning role
# Setting: {{ tuning_transparent_hugepages }}

[Unit]
Description=Configure Transparent Huge Pages ({{ tuning_transparent_hugepages }})
Documentation=https://www.kernel.org/doc/html/latest/admin-guide/mm/transhuge.html
DefaultDependencies=no
After=sysinit.target local-fs.target
Before=basic.target

[Service]
Type=oneshot
RemainAfterExit=yes

# Set THP enabled setting
ExecStart=/bin/sh -c 'echo {{ tuning_transparent_hugepages }} | tee /sys/kernel/mm/transparent_hugepage/enabled > /dev/null'

# Set THP defrag setting to same value
ExecStart=/bin/sh -c 'echo {{ tuning_transparent_hugepages }} | tee /sys/kernel/mm/transparent_hugepage/defrag > /dev/null'

# Verify the setting was applied
ExecStartPost=/bin/sh -c 'grep -q "\\[{{ tuning_transparent_hugepages }}\\]" /sys/kernel/mm/transparent_hugepage/enabled'

[Install]
WantedBy=basic.target
