# roles/tuning/templates/etc/systemd/system/network-tuning.service.j2
# Network Performance Tuning Service
# Generated by Ansible system tuning role
# Date: {{ ansible_date_time.iso8601 }}

[Unit]
Description=Network Performance Tuning
Documentation=man:ethtool(8)
After=network.target
Wants=network.target

[Service]
Type=oneshot
RemainAfterExit=yes

{% if valid_network_buffers is defined %}
{%   for buffer in valid_network_buffers %}
# Configure ring buffers for {{ buffer.interface }} (detected: {{ buffer.speed | default('unknown') }} speed)
ExecStart=/usr/sbin/ethtool -G {{ buffer.interface }} rx {{ buffer.rx_buffer }} tx {{ buffer.tx_buffer }}
ExecStartPost=/bin/sh -c '/usr/sbin/ethtool -g {{ buffer.interface }} | grep -q "RX:.*{{ buffer.rx_buffer }}"'
{%   endfor %}
{% elif tuning_network_ring_buffers is defined %}
{%   for buffer in tuning_network_ring_buffers %}
# Configure ring buffers for {{ buffer.interface }} (manual configuration)
ExecStart=/usr/sbin/ethtool -G {{ buffer.interface }} rx {{ buffer.rx_buffer }} tx {{ buffer.tx_buffer }}
ExecStartPost=/bin/sh -c '/usr/sbin/ethtool -g {{ buffer.interface }} | grep -q "RX:.*{{ buffer.rx_buffer }}"'
{%   endfor %}
{% else %}
# No network ring buffer configuration specified
ExecStart=/bin/true
{% endif %}

# Restart policy for network changes
Restart=no

[Install]
WantedBy=multi-user.target
