# roles/motd/templates/etc/update-motd.d/00-header.j2
#!/bin/sh
# {{ ansible_managed }}
{% if motd_header is defined %}
echo '{{ motd_header }}'
{% endif %}
